# conf/gunicorn.conf
description "Run Django instance"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec ./root/berkerpeksag/manage.py run_gunicorn localhost:8000