# conf/gunicorn.conf
description "Run Django instance"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec ./berkerpeksag/manage.py run_gunicorn 127.0.0.1:8000